# cm/agent/udev_monitor/99-nvme-hotplug.rules
ACTION=="add|remove", SUBSYSTEM=="block", RUN+="/usr/dragonfly/cm/agent/udev_monitor/SIGUSR1_kvcli.sh"
