# ansible/roles/ufm_etcd/templates/ufm_etcd.service.j2
[Unit]
Description=etcd on {{ inventory_hostname_short }} (This file is generated by Ansible)
Documentation=https://github.com/etcd
Conflicts=etcd.service
After=network.target

[Service]
Type=notify
Restart=on-failure
RestartSec=5
LimitNOFILE=40000
TimeoutStartSec=0

ExecStart=/usr/bin/etcd \
    --name {{ inventory_hostname_short }} \
    --data-dir /var/lib/etcd \
    --listen-client-urls http://{{ ansible_default_ipv4.address }}:{{ ufm_etcd_client_port }},http://127.0.0.1:{{ ufm_etcd_client_port }} \
    --advertise-client-urls http://{{ ansible_default_ipv4.address }}:{{ ufm_etcd_client_port }} \
    --listen-peer-urls http://{{ ansible_default_ipv4.address }}:{{ ufm_etcd_peer_port }} \
    --initial-advertise-peer-urls http://{{ ansible_default_ipv4.address }}:{{ ufm_etcd_peer_port }} \
    --initial-cluster {% for host in ansible_play_hosts -%}
            {%- if not loop.first -%}
                ,
            {%- endif -%}
            {{ hostvars[host].inventory_hostname_short }}=http://{{ hostvars[host].ansible_default_ipv4.address }}:{{ ufm_etcd_peer_port -}}
    {% endfor %}
      \
    --initial-cluster-token {{ ufm_etcd_token }} \
    --initial-cluster-state new

[Install]
WantedBy=multi-user.target
